# demo-2/tailon-supervisor.conf
[program:tailon]
command=/usr/local/bin/tailon -c /etc/tailon/config.yaml -d
stdout_logfile=/var/log/tailon-log.log
autostart=true
autorestart=true
directory=/etc/tailon
numprocs=1
user=ubuntu
